# dist/debian-init-script
#!/bin/sh

### BEGIN INIT INFO
# Provides:          nickserver
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: nickserver initscript
# Description:       Controls nickserver daemon (see https://github.com/leapcode/nickserver)
### END INIT INFO

PATH=/sbin:/usr/sbin:/bin:/usr/bin
NAME=nickserver
DAEMON=`which $NAME`

# exit if the daemon doesn't exist
[ -x "$DAEMON" ] || exit 0

. /lib/init/vars.sh
. /lib/lsb/init-functions

if [ "$VERBOSE" != no ]; then
  OPTIONS="--verbose"
else
  OPTIONS=""
fi

case "$1" in
  start)
    $DAEMON $OPTIONS start
    exit $?
    ;;
  stop)
    $DAEMON $OPTIONS stop
    exit $?
    ;;
  restart)
    $DAEMON $OPTIONS restart
    exit $?
    ;;
  status)
    $DAEMON $OPTIONS status
    exit $?
    ;;
  *)
    echo "Usage: /etc/init.d/$NAME {start|stop|restart|status}"
    exit 1
esac

exit 0